# lede/etc/uci-defaults/9999_enable_transitd
#!/usr/bin/env sh

set -e

/etc/init.d/transitd enable
/etc/init.d/transitd-redirect enable

exit 0
