# lede/etc/uci-defaults/9999_enable_transitd
#!/usr/bin/env sh

set -e

/etc/init.d/transitd enable && /etc/init.d/transitd start
/etc/init.d/transitd-redirect enable && /etc/init.d/transitd-redirect start
return $?
